# servers/supervisor/conf.d/php-fpm.ini
[program:php-fpm]
command=php-fpm -F
directory=/usr/local/sbin/
autostart=true
autorestart=true
startsecs=3
startretries=3
stdout_logfile=/var/log/php/supervisor/php-fpm.out.log
stderr_logfile=/var/log/php/supervisor/php-fpm.err.log
stdout_logfile_maxbytes=2MB
stderr_logfile_maxbytes=2MB
user=root
priority=999
numprocs=1
process_name=%(program_name)s_%(process_num)02d
stopasgroup=true
killasgroup=true